# screengrabber.timer
[Unit]
Description=Screen Grabber timer

[Timer]
OnCalendar=*-*-* *:*/10:00

[Install]
WantedBy=screengrabber.service
